# borg_no-stow/borgbackup-system.service
[Unit]
Description=Borg System Backup

[Service]
Type=oneshot
ExecStart=/root/bin/borgbackup-system.sh
